# deployment/systemd/eap-supercollider.service
[Unit]
Description=Electroacoustic Playground SuperCollider
After=eap-jack.service eap-subsequence.service
Requires=eap-jack.service
Wants=eap-subsequence.service

[Service]
Type=simple
User=we
Group=audio
WorkingDirectory=/opt/electroacoustic-playground
Environment=JACK_NO_AUDIO_RESERVATION=1
Environment=QT_QPA_PLATFORM=offscreen
# Port-name lists only — never channel counts (e.g. "4" loops SC outs into SC ins).
Environment=SC_JACK_DEFAULT_INPUTS=
Environment=SC_JACK_DEFAULT_OUTPUTS=system:playback_1,system:playback_2
EnvironmentFile=-/etc/default/eap-vcv
Environment=EAP_VCV_MIDI_DEVICE=14:0
LimitRTPRIO=95
LimitMEMLOCK=infinity
ExecStart=/bin/bash /opt/electroacoustic-playground/deployment/scripts/eap-supercollider-watchdog.sh
Restart=on-failure
RestartSec=2

[Install]
WantedBy=multi-user.target
