# deployment/systemd/eap-launchpad.service
[Unit]
Description=Electroacoustic Playground Launchpad controller
After=sound.target eap-supercollider.service
Requires=eap-supercollider.service
Wants=sound.target

[Service]
Type=simple
User=we
Group=audio
Environment=PYTHONUNBUFFERED=1
ExecStart=/usr/local/bin/eap-launchpad
Restart=on-failure
RestartSec=2

[Install]
WantedBy=multi-user.target
